# smart-hddtemp-daemon.service
[Unit]
Description=Smart hddtemp Daedmon

[Service]
ExecStart=/usr/bin/smart-hddtemp-daemon

[Install]
WantedBy=multi-user.target
